# overlay/etc/init.d/usb-tethering
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/lxc/files/lxc.initd.2,v 1.5 2012/07/21 05:07:15 flameeyes Exp $

description="USB tethering service to get debug USB network on Android devices"

command="/usr/bin/usb-tethering"

depend() {
        after bootmisc ssh
}
